# execucao_mensal.cron
00 12 10 * * /projects/das-scraper/run-das.sh >> /projects/das-scraper/logs/scraper.log 2>&1